# tools/deployment/com.allternitchitech.dak-runner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.allternit.dak-runner</string>

    <key>LabelDescription</key>
    <string>Allternit DAK Runner Daemon - Continuous Agent Execution</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Users/macbook/Desktop/allternit-workspace/allternit/1-kernel/agent-systems/allternit-dak-runner/dist/daemon.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/macbook/Desktop/allternit-workspace/allternit</string>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <true/>
        <key>Crashed</key>
        <true/>
        <key>NetworkState</key>
        <true/>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>Allternit_RAILS_CLI</key>
        <string>/Users/macbook/Desktop/allternit-workspace/allternit/target/release/allternit-rails</string>
        <key>Allternit_PROJECT_PATH</key>
        <string>/Users/macbook/Desktop/allternit-workspace/allternit</string>
        <key>Allternit_OUTPUT_DIR</key>
        <string>/Users/macbook/Desktop/allternit-workspace/allternit/.allternit/runner</string>
        <key>Allternit_MAX_CONCURRENT</key>
        <string>3</string>
        <key>Allternit_BUDGET_MAX</key>
        <string>90</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/macbook/Library/Logs/allternit/dak-runner.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/macbook/Library/Logs/allternit/dak-runner.err.log</string>

    <key>Nice</key>
    <integer>10</integer>

    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>512</integer>
        <key>MemoryLock</key>
        <integer>1073741824</integer>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>AbandonProcessGroup</key>
    <false/>
</dict>
</plist>
