# cmd/gizzi-code/cli-package/install/com.a2r.gizzi.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.a2r.gizzi</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/USERNAME/.local/bin/gizzi</string>
        <string>daemon</string>
        <string>start</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/USERNAME/Library/Logs/gizzi/output.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/USERNAME/Library/Logs/gizzi/error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/USERNAME/.local/bin</string>
        <key>GIZZI_LOG_LEVEL</key>
        <string>info</string>
    </dict>
</dict>
</plist>
