# infrastructure/chrome-stream/supervisor/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
rpcinterface=supervisor.supervisor_rpc_interface:SupervisorRPCInterface
http_port=9001
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info

[inet_http_server]
port=127.0.0.1:9001

[program:xvfb]
command=/usr/bin/Xvfb :99 -screen 0 %(ENV_A2R_RESOLUTION)sx24 -ac
autostart=true
autorestart=true
priority=10
user=root
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb_error.log
stdout_logfile_maxbytes=5MB
stderr_logfile_maxbytes=5MB
environment=DISPLAY=":99"

[program:dbus-daemon]
command=/usr/bin/dbus-daemon --system --fork
autostart=true
autorestart=true
priority=5
user=root
stdout_logfile=/var/log/supervisor/dbus.log
stderr_logfile=/var/log/supervisor/dbus_error.log
environment=DBUS_SESSION_BUS_ADDRESS="unix:path=/tmp/dbus-session"

[program:pulseaudio]
command=/usr/bin/pulseaudio --start --exit-idle-time=-1 --fail
autostart=true
autorestart=true
priority=6
user=chrome
stdout_logfile=/var/log/supervisor/pulseaudio.log
stderr_logfile=/var/log/supervisor/pulseaudio_error.log
environment=DISPLAY=":99"

[program:chrome]
command=/usr/bin/google-chrome-stable %(ENV_A2R_CHROME_FLAGS)s
autostart=true
autorestart=true
priority=20
user=chrome
stdout_logfile=/var/log/supervisor/chrome.log
stderr_logfile=/var/log/supervisor/chrome_error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=DISPLAY=":99",DBUS_SESSION_BUS_ADDRESS="unix:path=/tmp/dbus-session"
depends_on=xvfb,dbus-daemon
startsecs=3
startretries=3
stopsignal=TERM
stopwaitsecs=10

[program:selkies]
command=/usr/local/bin/selkies-gstreamer --addr=0.0.0.0 --port=8080 --encoder=%(ENV_SELKIES_ENCODER)s --enable_resize=true
autostart=true
autorestart=true
priority=30
user=root
stdout_logfile=/var/log/supervisor/selkies.log
stderr_logfile=/var/log/supervisor/selkies_error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=DISPLAY=":99",GI_TYPELIB_PATH="/usr/lib/x86_64-linux-gnu/girepository-1.0",LD_LIBRARY_PATH="/usr/lib/x86_64-linux-gnu:/usr/lib",GST_PLUGIN_PATH="/usr/lib/x86_64-linux-gnu/gstreamer-1.0",GST_PLUGIN_SYSTEM_PATH="/usr/lib/x86_64-linux-gnu/gstreamer-1.0",PYTHONPATH="/usr/lib/python3/dist-packages",PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
depends_on=chrome
startsecs=2
startretries=3

[program:sidecar]
command=/usr/local/bin/uvicorn main:app --host 0.0.0.0 --port 8081 --workers 1
autostart=true
autorestart=true
priority=40
user=root
directory=/opt/a2r/sidecar
stdout_logfile=/var/log/supervisor/sidecar.log
stderr_logfile=/var/log/supervisor/sidecar_error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
startsecs=2
startretries=3
