# api/cloud/allternit-node/deploy/services/allternit-node.service
[Unit]
Description=Allternit Node Agent - Connect your infrastructure to Allternit Control Plane
Documentation=https://docs.allternit.io
After=network.target network-online.target docker.service
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/allternit/node.env
ExecStart=/opt/allternit/bin/allternit-node
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=always
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=allternit-node

# User/Group
User=root
Group=root

# Working directory
WorkingDirectory=/opt/allternit

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true

# Allow Docker socket access
ReadWritePaths=/var/run/docker.sock

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=2G
CPUQuota=200%

# Network
IPAddressDeny=any
IPAddressAllow=localhost
IPForward=true

# System call filtering
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources

# Capabilities
CapabilityBoundingSet=CAP_NET_BIND_SERVICE CAP_SYS_ADMIN
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
